# data/hpm-mux.service
[Unit]
Description=OMEN Command Center — MUX (GPU Switch) Service
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/libexec/hp-manager/services/mux_service.py
Restart=on-failure
RestartSec=3

# Security hardening
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
LockPersonality=true
RestrictSUIDSGID=true

# Needs access to GPU switching tools and sysfs
ReadWritePaths=/sys/devices/platform/hp-wmi /etc/hp-manager

[Install]
WantedBy=multi-user.target
